# host/systemd/shared-store.service
[Unit]
Description=Auto-unlock and mount the shared LUKS store volume
Documentation=file:///opt/raph-server-installer/docs/design.md
# Both the blob source and the mountpoint parent must exist before we run.
RequiresMountsFor=/srv/store/data /srv/store/mnt
# Soft preconditions — the unit becomes a no-op (logged, not failed) on
# any host where the shared volume hasn't been provisioned yet (e.g. a
# partial bootstrap, or a recovery rebuild where create-shared-volume.sh
# hasn't been re-run yet).
ConditionPathExists=/etc/luks/_shared.key
ConditionPathExists=/srv/store/data/_shared.img
# Order before docker so the mount exists before any container that
# bind-mounts /srv/store/mnt (cloud, in particular) is brought up.
After=local-fs.target
Before=docker.service
DefaultDependencies=yes

[Service]
Type=oneshot
RemainAfterExit=true

# Open the LUKS blob with the host-side keyfile (no passphrase prompt).
# Same crypto profile as the per-user store-mount@ unit; only the unlock
# secret differs. Mapper name `store_shared` matches what
# create-shared-volume.sh formats the volume against.
ExecStart=/sbin/cryptsetup open --type luks2 --key-file=/etc/luks/_shared.key /srv/store/data/_shared.img store_shared
ExecStart=/bin/mkdir -p /srv/store/mnt/_shared
ExecStart=/bin/mount /dev/mapper/store_shared /srv/store/mnt/_shared

# Clean teardown. `-` prefix on umount lets close still run if the mount
# was already gone. We deliberately do NOT --deferred-close: the unit is
# expected to be the only consumer of this mapper.
ExecStop=-/bin/umount /srv/store/mnt/_shared
ExecStop=/sbin/cryptsetup close store_shared

[Install]
WantedBy=multi-user.target
